# hid-rebind.rules
SUBSYSTEM=="hid", ACTION=="add", ENV{HID_ID}=="0003:0000256C:*", \
	RUN+="/sbin/hid-rebind"
